# ospf-latensy-cost-autocalc - linux/ospf-latensy-cost-autocalc.service
[Unit]
Description=Monitor OSPF neighbors and adjust cost

[Service]
ExecStart=%h/ospf-latensy-cost-autocalc.py